# init-scripts/abrt-oops.service
[Unit]
Description=ABRT kernel log watcher
After=abrtd.service
Requires=abrtd.service

[Service]
ExecStart=/usr/bin/abrt-dump-oops -rwxD /var/log/messages

[Install]
WantedBy=multi-user.target
